# sensord-1.1.17/sensord.manifest
<manifest>
	<define>
		<domain name="sensord" />
		<provide>
			<label name="sensord::server" />
			<label name="sensord::bio" />
		</provide>
	</define>
	<request>
		<domain name="sensord" />
	</request>
	<assign>
		<filesystem path="/etc/rc.d/init.d/sfsvc" label="_" exec_label="none" />
	</assign>
</manifest>
