# sensord-1.1.17/packaging/sensord.socket
[Unit]
Description=Sensor daemon socket

[Socket]
ListenStream=/tmp/sf_command_socket
SocketMode=0777
PassCredentials=yes
Accept=false
SmackLabelIPIn=*
SmackLabelIPOut=@
